# files/rc.local
/usr/sbin/ntpdate 192.168.105.41

touch /var/lock/subsys/local
route add -net 192.168.0.0 netmask 255.255.0.0 gw 10.64.7.254
route add -net 10.0.0.0 netmask 255.0.0.0 gw 10.64.7.254
/usr/local/lvs-dr-rs.sh start
sleep 1
/usr/local/nginx/sbin/nginx -c /usr/local/nginx/conf/nginx.conf
/opt/memcached/bin/memcached -d -uroot -m 1024 -l 127.0.0.1 -p 11211
cd /root/cdn/ && /root/cdn/safe_ccws.sh &
